# lib/capistrano/parcel/supervisor/service.conf.erb
[program:<%= fetch(:application) %>]
command=<%= fetch(:uwsgi_bin) %> --ini=/etc/uwsgi/<%= fetch(:application) %>.ini
process_name=<%= fetch(:application) %>
numprocs=<%= fetch(:supervisor_numprocs) %>
directory=/tmp
umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=<%= fetch(:supervisor_user) %>
redirect_stderr=false
stdout_logfile=<%= fetch(:log_dir) %>/<%= fetch(:application) %>.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
serverurl=AUTO